# pkg/sevbesau/files-base/etc/profile.d/sevbesau-base.sh
# Managed by sevbesau-base package
# https://github.com/sevbesau/arch-repo
# Initialize variables for login shells only


export EDITOR='nvim'
export PATH=$HOME/.local/bin:$PATH
export TERM=xterm-256color

# Node
export PATH=$HOME/.local/share/npm/bin:$HOME/.yarn/bin:$PATH
