# pkg/sevbesau/sevbesau-desktop.install
post_install() {
    post_upgrade
}

post_upgrade() {
    # Start and enable the bluetooth service
    systemctl enable --now bluetooth.service

    # Reload udev rules
    systemctl restart systemd-udevd
}