# pkg/sevbesau/sevbesau-desktop.install
post_install() {
    post_upgrade
}

post_upgrade() {
    # start and enable pipewire-pulse.sevice
    sleep 1
    
    systemctl enable --now bluetooth.service

}